# vp-bot.service
[Unit]
Description=Discord VP Tracker Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USERNAME
WorkingDirectory=/home/YOUR_LINUX_USERNAME/vp-bot
ExecStart=/home/YOUR_LINUX_USERNAME/vp-bot/venv/bin/python3 bot.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
